# systemd/sync-ssh-ecs.service
[Unit]
Description=Sync ECS instance IPs to SSH config

[Service]
Type=oneshot
WorkingDirectory=/home/martin/Documents/Programming/MESH/tools
ExecStart=/home/martin/.local/bin/uv run /home/martin/Documents/Programming/MESH/tools/sync_ssh.py --profile mesh
